# install/deployment/creaturebox.service
# deployment/creaturebox.service
[Unit]
Description=CreatureBox Web Interface
After=network.target
Wants=redis-server.service

[Service]
# Use environment variables to make paths configurable
Environment="CREATUREBOX_HOME=/opt/creaturebox"
Environment="CREATUREBOX_VENV=/opt/creaturebox-venv"
Environment="CREATUREBOX_LOG_DIR=/var/log/creaturebox"
Environment="CREATUREBOX_PORT=5000"
Environment="CREATUREBOX_LOG_LEVEL=INFO"

# Override these using systemd override files rather than editing this file
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="ENABLE_CAMERA_STREAM=true"
Environment="ENABLE_BACKGROUND_JOBS=true"
Environment="ENABLE_RATE_LIMITING=true"

# Path configuration
WorkingDirectory=${CREATUREBOX_HOME}

# Use proper user
User=creaturebox
Group=creaturebox

# Execution command using Gunicorn
ExecStart=${CREATUREBOX_VENV}/bin/gunicorn \
    -c ${CREATUREBOX_HOME}/deployment/gunicorn.conf.py \
    --chdir ${CREATUREBOX_HOME} \
    --bind 127.0.0.1:${CREATUREBOX_PORT} \
    'src.web.app:create_app()'

# Restart policy
Restart=on-failure
RestartSec=5

# Ensure proper shutdown
KillSignal=SIGTERM
TimeoutStopSec=20

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target